# incident-log.service
[Unit]
Description=Incident Log System
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/incident-log
Environment="PATH=/opt/incident-log/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_APP=app.py"
ExecStart=/opt/incident-log/venv/bin/python app.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/incident-log/app.log
StandardError=append:/var/log/incident-log/error.log

[Install]
WantedBy=multi-user.target
